# deployment/systemd/safebox-web.service
[Unit]
Description=SafeBox Local Web UI
After=network-online.target
Wants=network-online.target
RequiresMountsFor=/opt

[Service]
Type=simple
User=__SERVICE_USER__
WorkingDirectory=__PROJECT_DIR__
ExecStart=__VENV_DIR__/bin/python -m web.app
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
